# ansible/templates/validator-launcher.service.j2
[Unit]
Description=Validator Launcher Service
After=network.target
Requires=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStartPre=/usr/local/bin/start-dstack-services.sh
ExecStart=/usr/local/bin/validator-launcher run
ExecStopPost=/usr/local/bin/stop-dstack-services.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="RUST_LOG=info"
Environment="VMM_URL=http://localhost:10300"
Environment="DSTACK_DIR=/home/ubuntu/meta-dstack/build"
Environment="LOG_DIR=/var/log/dstack"

[Install]
WantedBy=multi-user.target

